# debian/extra/rules-deepin/81-enable-wol.rules
ACTION=="add", SUBSYSTEM=="net", NAME=="en*", PROGRAM=="/usr/bin/test -x /sbin/ethtool", RUN+="/sbin/ethtool -s $name wol g"
